# deploy/systemd/things-sentiment-poller.service
# Systemd service for OpenClaw Things Sentiment Polling Daemon
# Installation:
#   1. Edit this file and set WORKDIR to your actual workspace path
#   2. Or use: env WORKDIR=/your/path ./scripts/install_systemd.sh
#   3. Then: systemctl --user daemon-reload && systemctl --user enable --now things-sentiment-poller

[Unit]
Description=OpenClaw Things Sentiment Polling Service
Documentation=https://github.com/faizanxg-AI/openclaw-things-sentiment
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# IMPORTANT: Set WorkingDirectory to your actual workspace path
WorkingDirectory=${WORKDIR:-/path/to/agent-bridge/workspace}
ExecStart=${WORKDIR:-/path/to/agent-bridge/workspace}/scripts/start_polling.sh
Environment="OPENCLAW_SESSION_KEY=${OPENCLAW_SESSION_KEY}"
Environment="PYTHONUNBUFFERED=1"
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=things-sentiment-poller

# Security hardening (uncomment if needed)
# PrivateTmp=true
# ProtectSystem=strict
# ProtectHome=read-only
# NoNewPrivileges=true

[Install]
WantedBy=default.target
